# devsupport/encoder/supervisor.d/enc-002.conf
[group:enc-2]
programs=audio-2,pad-2

[program:audio-2]
directory=/encoder
environment=
  SRC="http://stream.openbroadcast.ch/128.mp3",
  ENCODER_ID="2",
  AUDIOENC_ARGS="--bitrate 64 --rate 48000 --channels 2 --pad 32"
command=/usr/bin/liquidsoap ls/encoder.liq
killasgroup=true
stopsignal=KILL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:pad-2]
directory=/encoder
command=odr-padenc
  --sleep 60
  --output enc-2
  --dls /encoder/run/enc-1/dls.txt
  --dir /encoder/run/enc-1/sls
  -v
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
