# devsupport/encoder/supervisor.d/enc-004.conf
[group:enc-4]
programs=audio-4,pad-4

[program:audio-4]
directory=/encoder
environment=
  SRC="http://stream.openbroadcast.ch/128.mp3",
  ENCODER_ID="4",
  AUDIOENC_ARGS="--bitrate 80 --rate 48000 --channels 2 --pad 64 --ps"
command=/usr/bin/liquidsoap ls/encoder.liq
killasgroup=true
stopsignal=KILL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:pad-4]
directory=/encoder
command=odr-padenc
  --sleep 10
  --output enc-4
  --dls /encoder/run/enc-1/dls.txt
  --dir /encoder/run/enc-1/sls
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true