# devsupport/encoder/supervisor.d/enc-003.conf
[group:enc-3]
programs=audio-3,pad-3

[program:audio-3]
directory=/encoder
environment=
  SRC="http://stream.openbroadcast.ch/128.mp3",
  ENCODER_ID="3",
  AUDIOENC_ARGS="--bitrate 72 --rate 48000 --channels 2 --pad 64"
command=/usr/bin/liquidsoap ls/encoder.liq
killasgroup=true
stopsignal=KILL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:pad-3]
directory=/encoder
command=odr-padenc
  --sleep 10
  --output enc-3
  --dls /encoder/run/enc-1/dls.txt
  --dir /encoder/run/enc-1/sls
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
