# devsupport/encoder/supervisor.d/enc-001.conf
[group:enc-1]
programs=audio-1,pad-1

[program:audio-1]
directory=/encoder
environment=
  SRC="http://stream.openbroadcast.ch/128.mp3",
  ENCODER_ID="1",
  AUDIOENC_ARGS="--bitrate 48 --rate 48000 --channels 2 --pad 16"
command=/usr/bin/liquidsoap ls/encoder.liq
killasgroup=true
stopsignal=KILL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:pad-1]
directory=/encoder
command=odr-padenc
  --sleep 60
  --output enc-1
  --dls /encoder/run/enc-1/dls.txt
  --dir /encoder/run/enc-1/sls
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
